# build/linux-server/crynux-node.service
[Unit]
Description=Start Crynux Node on System Startup

[Service]
WorkingDirectory=/usr/local/crynux-node
ExecStart=/usr/local/crynux-node/crynux-node
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
